# modules/dirtyfix/files/rc.local
#!/bin/sh -e

chmod 666 /dev/mali0

exit 0
